# .config/autostart/guake.desktop
 [Desktop Entry]
 Type=Application
 Terminal=false
 Name=Guake
 Exec=guake
 Icon=
 Comment=Run guake
 Categories=System;